# celery.service
[Unit]
Description=0123.ru Celery Worker
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/var/www/filehost
Environment="PATH=/var/www/filehost/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=filehost.settings_prod"
Environment="PYTHONPATH=/var/www/filehost"
ExecStart=/var/www/filehost/venv/bin/celery multi start worker1 -A filehost --loglevel=info --pidfile=/var/run/celery/worker.pid --logfile=/var/log/celery/worker.log
ExecStop=/var/www/filehost/venv/bin/celery multi stopwait worker1 --pidfile=/var/run/celery/worker.pid
ExecReload=/var/www/filehost/venv/bin/celery multi restart worker1 -A filehost --loglevel=info --pidfile=/var/run/celery/worker.pid --logfile=/var/log/celery/worker.log
Restart=always
RestartSec=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/filehost/media /var/www/filehost/logs /var/run/celery /var/log/celery
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=celery

[Install]
WantedBy=multi-user.target
